# deploy/host/systemd/helpdesk-auditd.service
[Unit]
Description=aiHelpDesk Audit Daemon
Documentation=https://github.com/borisdali/helpdesk
After=network.target
# Start before agents so audit logging is available from the first tool call.

[Service]
Type=simple
User=helpdesk
Group=helpdesk
WorkingDirectory=/opt/helpdesk
EnvironmentFile=/opt/helpdesk/.env
ExecStart=/opt/helpdesk/auditd
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helpdesk-auditd

# Harden the process: read-only filesystem except for the audit database.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/helpdesk /var/lib/helpdesk

[Install]
WantedBy=helpdesk.target
