# deploy/host/systemd/helpdesk-gateway.service
[Unit]
Description=aiHelpDesk Gateway
Documentation=https://github.com/borisdali/helpdesk
After=network.target helpdesk-auditd.service helpdesk-database-agent.service helpdesk-k8s-agent.service helpdesk-sysadmin-agent.service helpdesk-incident-agent.service
Wants=helpdesk-auditd.service helpdesk-database-agent.service helpdesk-k8s-agent.service helpdesk-sysadmin-agent.service helpdesk-incident-agent.service

[Service]
Type=simple
User=helpdesk
Group=helpdesk
WorkingDirectory=/opt/helpdesk
EnvironmentFile=/opt/helpdesk/.env
ExecStart=/opt/helpdesk/gateway
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=helpdesk-gateway

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/helpdesk /var/lib/helpdesk

[Install]
WantedBy=helpdesk.target
